# debian/flocker-buildslave.upstart
# Upstart service for buildslave.
start on (local-filesystems and net-device-up IFACE!=lo)
exec /usr/bin/twistd -ny /srv/flocker-buildslave/slave/buildbot.tac
setuid flocker-buildslave
chdir /srv/flocker-buildslave/slave
