# debian/flocker-buildslave.upstart
# Upstart service for buildslave.
start on (local-filesystems and net-device-up IFACE!=lo)
exec /usr/bin/twistd -ny /srv/buildslave/slave/buildbot.tac
setuid buildslave
chdir /srv/buildslave/slave
